# src/flower/DEBIAN/prerm
#!/bin/sh

systemctl disable flower

systemctl stop flower

exit 0
